# default.cron
* * * * * /usr/local/bin/php -f /app/src/run.php
