# default.cron
* * * * * /usr/local/bin/php -q /app/src/run.php --run
